# image_processing/rc.local
/home/ubuntu/Documents/FIRST_Image_Processing/image_processing/startup_processing.sh
